# deploy/warp-dashboard-device3.service
[Unit]
Description=WARP Dashboard Device 3 (Vite)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/warp/app/dashboard
ExecStart=/usr/bin/npm run preview -- --host 0.0.0.0 --port 5173
Restart=always
EnvironmentFile=/etc/warp/env
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
