# deploy/warp-ui-device3.service
[Unit]
Description=WARP UI Device 3 (Next.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/warp/app/ui/warp-console
ExecStart=/usr/bin/npm run start
Restart=always
EnvironmentFile=/etc/warp/env
Environment=NODE_ENV=production
Environment=PORT=3003

[Install]
WantedBy=multi-user.target
