# deploy/warp-device3.service
[Unit]
Description=WARP Device 3 Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/warp/app/devices/device3
ExecStart=/opt/warp/venv/bin/python -m src.main --config /etc/warp/device3.yaml
Restart=always
EnvironmentFile=/etc/warp/env

[Install]
WantedBy=multi-user.target
